# tt-device-queue.service
[Unit]
Description=Tenstorrent device queue - durable serialized hardware access
After=network.target

[Service]
Type=simple
WorkingDirectory=%h
ExecStart=%h/tenstorrent/tt-device-queue/.venv/bin/python3 %h/tenstorrent/tt-device-queue/server.py
Restart=on-failure
RestartSec=2
TimeoutStopSec=20
KillMode=control-group
OOMPolicy=continue
LimitNOFILE=65536
NoNewPrivileges=yes
Environment=PATH=/usr/local/bin:/usr/bin:/bin:%h/.local/bin
Environment=PYTHONUNBUFFERED=1
Environment=TT_DEVICE_LOG_DIR=%h/tenstorrent/tt-device-queue/logs-v2

[Install]
WantedBy=default.target
